# .config/systemd/user/doorlocker-controller@.service
[Unit]
Description=auto door locker contoll service service

[Service]
ExecStart=%h/doorlocker/door.bash
StandardInput=socket
StandardOutput=socket

