# .config/systemd/user/doorlocker-controller.socket
[Unit]
Description=auto door locker contoll service socket

[Socket]
ListenStream=127.0.0.1:2048
Accept=true

[Install]
WantedBy=sockets.target

